# root/etc/s6-overlay/s6-rc.d/init-keygen-perms/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# permissions
aetherxown -R dockerx:dockerx \
    /config/keys
