# root/etc/s6-overlay/s6-rc.d/init-custom-files/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# #
#   define > colors
#
#   Use the color table at:
#       - https://gist.github.com/fnky/458719343aabd01cfb17a3a4f7296797
# #

declare -A c=(
    [end]=$'\e[0m'
    [white]=$'\e[97m'
    [bold]=$'\e[1m'
    [dim]=$'\e[2m'
    [underline]=$'\e[4m'
    [strike]=$'\e[9m'
    [blink]=$'\e[5m'
    [inverted]=$'\e[7m'
    [hidden]=$'\e[8m'
    [black]=$'\e[0;30m'
    [redl]=$'\e[0;91m'
    [redd]=$'\e[0;31m'
    [magental]=$'\e[0;95m'
    [magentad]=$'\e[0;35mm'
    [bluel]=$'\e[0;94m'
    [blued]=$'\e[0;34m'
    [cyanl]=$'\e[0;96m'
    [cyand]=$'\e[0;36m'
    [greenl]=$'\e[0;92m'
    [greend]=$'\e[0;32m'
    [yellowl]=$'\e[0;93m'
    [yellowd]=$'\e[0;33m'
    [greyl]=$'\e[0;37m'
    [greyd]=$'\e[0;90m'
    [navy]=$'\e[38;5;62m'
    [olive]=$'\e[38;5;144m'
    [peach]=$'\e[38;5;210m'
)

# #
#   unicode for emojis
#       https://apps.timwhitlock.info/emoji/tables/unicode
# #

declare -A icon=(
    ["symbolic link"]=$'\xF0\x9F\x94\x97' # ðŸ”—
    ["regular file"]=$'\xF0\x9F\x93\x84' # ðŸ“„
    ["directory"]=$'\xF0\x9F\x93\x81' # ðŸ“
    ["regular empty file"]=$'\xe2\xad\x95' # â­•
    ["log"]=$'\xF0\x9F\x93\x9C' # ðŸ“œ
    ["1"]=$'\xF0\x9F\x93\x9C' # ðŸ“œ
    ["2"]=$'\xF0\x9F\x93\x9C' # ðŸ“œ
    ["3"]=$'\xF0\x9F\x93\x9C' # ðŸ“œ
    ["4"]=$'\xF0\x9F\x93\x9C' # ðŸ“œ
    ["5"]=$'\xF0\x9F\x93\x9C' # ðŸ“œ
    ["pem"]=$'\xF0\x9F\x94\x92' # ðŸ”‘
    ["pub"]=$'\xF0\x9F\x94\x91' # ðŸ”’
    ["pfx"]=$'\xF0\x9F\x94\x92' # ðŸ”‘
    ["p12"]=$'\xF0\x9F\x94\x92' # ðŸ”‘
    ["key"]=$'\xF0\x9F\x94\x91' # ðŸ”’
    ["crt"]=$'\xF0\x9F\xAA\xAA ' # ðŸªª
    ["gz"]=$'\xF0\x9F\x93\xA6' # ðŸ“¦
    ["zip"]=$'\xF0\x9F\x93\xA6' # ðŸ“¦
    ["gzip"]=$'\xF0\x9F\x93\xA6' # ðŸ“¦
    ["deb"]=$'\xF0\x9F\x93\xA6' # ðŸ“¦
    ["sh"]=$'\xF0\x9F\x97\x94' # ðŸ—”
)

# Directories
SCRIPTS_DIR="/custom-cont-init.d"

# Make sure custom init directory exists and has files in it
if [[ -e "${SCRIPTS_DIR}" ]] && [[ -n "$(/bin/ls -A ${SCRIPTS_DIR} 2>/dev/null)" ]]; then
    printf '%-29s %-65s\n' "  ${c[bluel]}Loader${c[end]}" "${c[end]}Loading any found plugins${c[end]}"
    for SCRIPT in "${SCRIPTS_DIR}"/*; do
        NAME="$(basename "${SCRIPT}")"
        if [[ -f "${SCRIPT}" ]]; then
            printf '%-29s %-65s\n' "  ${c[bluel]}Loader${c[end]}" "${c[end]}Executing${c[end]}"
            /bin/bash "${SCRIPT}"
            printf '%-29s %-65s\n' "  ${c[bluel]}Loader${c[end]}" "${c[end]}Successfully ran with code ${c[bluel]}[$?]${c[end]}"
        elif [[ ! -f "${SCRIPT}" ]]; then
            printf '%-29s %-65s\n' "  ${c[bluel]}Loader${c[end]}" "${c[end]}${c[bluel]}${NAME}${c[end]} is not a valid file${c[end]}"
        fi
    done
else
    printf '%-29s %-65s\n' "  ${c[bluel]}Loader${c[end]}" "${c[end]}No plugins found; skipping${c[end]}"
fi
